# systemd/storm-health-monitor.service
[Unit]
Description=STORM Health Monitor and Auto Failover
After=network-online.target storm-client.service storm-resolver-daemon.service
Wants=network-online.target storm-client.service storm-resolver-daemon.service

[Service]
Type=simple
WorkingDirectory=/opt/nexora-storm
ExecStart=/opt/nexora-storm/.venv/bin/python /opt/nexora-storm/storm_health_monitor.py --proxy-host 127.0.0.1 --proxy-port 1443 --target-host 1.1.1.1 --target-port 443 --targets-file /opt/nexora-storm/state/monitor_targets.txt --checks 6 --interval 0.5 --timeout 6 --success-threshold 0.50 --latency-threshold-ms 5000 --loop 30 --active-file /opt/nexora-storm/state/resolvers_active.txt --healthy-file /opt/nexora-storm/state/resolvers_healthy.txt --take 4 --fail-streak-trigger 3 --action-cooldown 180 --failover-rotate 1 --e2e-rank-min-fail-streak 6 --e2e-rank-cmd "/opt/nexora-storm/.venv/bin/python /opt/nexora-storm/storm_resolver_e2e_ranker.py" --state-json /opt/nexora-storm/state/health_monitor_state.json --report-json /opt/nexora-storm/state/health_monitor_report.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
