# systemd/storm-client.service
[Unit]
Description=STORM Inside Client (Auto Resolver Selection)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/nexora-storm
ExecStart=/opt/nexora-storm/run_storm_client_auto.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

